# draggy.service
[Unit]
Description=draggy

[Service]
Type=simple
ExecStart=draggy

[Install]
WantedBy=graphical.target
